# ansible/roles/scanner/templates/scanner-trigger.path.j2
# ANSIBLE MANAGED — scanner role
[Unit]
Description=Watch for scanner button trigger flag

[Path]
PathExists={{ scanner_trigger_flag }}
Unit=scanner-trigger.service

[Install]
WantedBy=multi-user.target
